# system/modules/home-pc-linux/fancontrol/root/etc/systemd/system/fancontrol.service.d/override.conf
[Service]
Restart=on-failure
RestartSec=1s
