# profiles/server/system/etc/systemd/system/qbittorrent-nox@.service.d/override.conf
[Service]
PrivateTmp=true
ReadOnlyPaths=/
TemporaryFileSystem=/var
ReadWritePaths=/home
BindPaths=/var/data/tmp
