# profiles/server/system/etc/udev/rules.d/69-hdparm.rules
ACTION=="add|change", SUBSYSTEM=="block", KERNEL=="sd[b-c]", ATTRS{queue/rotational}=="1", RUN+="/usr/bin/hdparm -S 24 /dev/%k"