# system/modules/server/files/root/etc/udev/rules.d/69-hdparm.rules
ACTION=="add|change", SUBSYSTEM=="block", KERNEL=="sd[b-f]", ATTRS{queue/rotational}=="1", RUN+="/usr/bin/hdparm -S 50 /dev/%k"
