# user/workstation/user/home/.config/autostart/sunshine.desktop
[Desktop Entry]
Name=Sunshine
GenericName=Desktop Streaming
Exec=/usr/bin/sunshine
Terminal=false
