# system/modules/lenovo-laptop/sleep/root/etc/systemd/system/resleep-hook.service
[Unit]
Description=Resleep hook - force sleep if lid is closed after wakeup
After=suspend.target

[Service]
Type=oneshot
ExecStart=/opt/resleep-hook.sh

[Install]
WantedBy=suspend.target
